# phoenix.service
[Unit]
Description=Phoenix service

[Service]
ExecStart=./start.sh

[Install]
WantedBy=multi-user.target